# scripts/openvpn-mgmt.service
[Unit]
Description=OpenVPN Management Service
After=network.target openvpn-server@server.service
Wants=openvpn-server@server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/openvpn-mgmt
ExecStart=/opt/openvpn-mgmt/openvpn-mgmt serve
Restart=always
RestartSec=10
Environment=DB_PATH=/opt/openvpn-mgmt/data/openvpn.db
Environment=EASYRSA_PATH=/etc/openvpn/easy-rsa
Environment=OPENVPN_PATH=/etc/openvpn
Environment=CLIENTS_DIR=/etc/openvpn/clients
Environment=JWT_SECRET_FILE=/opt/openvpn-mgmt/.jwt_secret

# Security
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/openvpn-mgmt/data /etc/openvpn/clients /etc/openvpn/easy-rsa/pki

[Install]
WantedBy=multi-user.target
